# icinga_network-services.timer
[Unit]
Description=Daily network scan
[Timer]
OnCalendar=daily
Persistent=true
[Install]
WantedBy=multi-user.target
