# icinga_network-services.service
[Unit]
Description=Network scanner
[Service]
Type=oneshot
ExecStart=/etc/icingaweb2/modules/network-services/scanner
SuccessExitStatus=1
